# configs/dev_celery_sms_worker.conf
[program:dev_sms_EGYPT_Mobinil_1000]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q sms_EGYPT_Mobinil_1000 --loglevel=info  -n sms_EGYPT_Mobinil_1000_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/sms_EGYPT_Mobinil_1000.log
stderr_logfile= /var/log/celery-billing/development/sms_EGYPT_Mobinil_1000.log
autostart= True
autorestart= False

[program:dev_sms_NIGERIA_Globacom_917]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q sms_NIGERIA_Globacom_917 --loglevel=info  -n sms_NIGERIA_Globacom_917_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_917.log
stderr_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_917.log
autostart= True
autorestart= False

[program:dev_sms_NIGERIA_Globacom_915]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q sms_NIGERIA_Globacom_915 --loglevel=info  -n sms_NIGERIA_Globacom_915_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_915.log
stderr_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_915.log
autostart= True
autorestart= False

[program:dev_sms_NIGERIA_Globacom_916]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q sms_NIGERIA_Globacom_916 --loglevel=info  -n sms_NIGERIA_Globacom_916_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_916.log
stderr_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_916.log
autostart= True
autorestart= False

[program:dev_sms_NIGERIA_Globacom_917]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q sms_NIGERIA_Globacom_917 --loglevel=info  -n sms_NIGERIA_Globacom_917_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_917.log
stderr_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_917.log
autostart= True
autorestart= False

[program:dev_sms_NIGERIA_Globacom_914]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q sms_NIGERIA_Globacom_914 --loglevel=info  -n sms_NIGERIA_Globacom_914_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_914.log
stderr_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_914.log
autostart= True
autorestart= False

[program:dev_sms_NIGERIA_Globacom_918]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q sms_NIGERIA_Globacom_918 --loglevel=info  -n sms_NIGERIA_Globacom_918_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_918.log
stderr_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_918.log
autostart= True
autorestart= False

[program:dev_sms_NIGERIA_Globacom_916]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q sms_NIGERIA_Globacom_916 --loglevel=info  -n sms_NIGERIA_Globacom_916_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_916.log
stderr_logfile= /var/log/celery-billing/development/sms_NIGERIA_Globacom_916.log
autostart= True
autorestart= False

