# configs/dev_celery_bill_worker.conf
[program:dev_bill_NIGERIA_Globacom_916]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q bill_NIGERIA_Globacom_916 --loglevel=info  -n bill_NIGERIA_Globacom_916_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_916.log
stderr_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_916.log
autostart= True
autorestart= False

[program:dev_bill_NIGERIA_Globacom_913]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q bill_NIGERIA_Globacom_913 --loglevel=info  -n bill_NIGERIA_Globacom_913_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_913.log
stderr_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_913.log
autostart= True
autorestart= False

[program:dev_bill_NIGERIA_Globacom_913]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q bill_NIGERIA_Globacom_913 --loglevel=info  -n bill_NIGERIA_Globacom_913_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_913.log
stderr_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_913.log
autostart= True
autorestart= False

[program:dev_bill_NIGERIA_Globacom_913]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q bill_NIGERIA_Globacom_913 --loglevel=info  -n bill_NIGERIA_Globacom_913_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_913.log
stderr_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_913.log
autostart= True
autorestart= False

[program:dev_bill_NIGERIA_Globacom_914]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q bill_NIGERIA_Globacom_914 --loglevel=info  -n bill_NIGERIA_Globacom_914_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_914.log
stderr_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_914.log
autostart= True
autorestart= False

[program:dev_bill_NIGERIA_Globacom_918]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q bill_NIGERIA_Globacom_918 --loglevel=info  -n bill_NIGERIA_Globacom_918_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_918.log
stderr_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_918.log
autostart= True
autorestart= True

[program:dev_bill_NIGERIA_Globacom_917]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q bill_NIGERIA_Globacom_917 --loglevel=info  -n bill_NIGERIA_Globacom_917_worker --concurrency=1  --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_917.log
stderr_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_917.log
autostart= True
autorestart= False

[program:dev_bill_NIGERIA_Globacom_916]
command=/home/ahussein/platform-celery/celery-billing/venv/bin/celery -A dev_celery_runner worker -Q bill_NIGERIA_Globacom_916 --autoscale=10,2 --loglevel=info  -n bill_NIGERIA_Globacom_916_worker --without-gossip --without-mingle
directory=/home/ahussein/platform-celery/celery-billing
user=billing
stdout_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_916.log
stderr_logfile= /var/log/celery-billing/development/bill_NIGERIA_Globacom_916.log
autostart= True
autorestart= False

